# systemd/user/offlineimap.timer
[Unit]
Description=Offlineimap Query Timer

[Timer]
OnUnitInactiveSec=30s
OnBootSec=1s
Unit=offlineimap.service

[Install]
WantedBy=default.target
