# systemd/user/offlineimap.service
[Unit]
Description=Offlineimap Service

[Service]
Type=simple
ExecStart=/usr/bin/offlineimap -o -q -u quiet
